# lg-service-openrc
#!/sbin/openrc-run

# Looking Glass OpenRC service

ouch /dev/shm/looking-glass && chown $VIRT_USER:kvm /dev/shm/looking-glass && chmod 660 /dev/shm/looking-glass"

exit
